# demo/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:code-server]
command=/usr/bin/code-server --bind-addr 0.0.0.0:8080 --auth none /home/coder/workspace/grafema
user=coder
environment=HOME="/home/coder"
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
